# deployment/crontab
0 15 * * * /Users/pli/miniconda3/envs/housing/bin/python /Users/pli/Documents/projects/housingfinder/main.py >> ~/cron.log 2>&1